# templates/default/ircd-services.init.erb
#!/bin/sh
### BEGIN INIT INFO
# Provides:          inspircd
# Required-Start:    $remote_fs $network $syslog $time
# Required-Stop:     $remote_fs $syslog
# Should-Start:      $local_fs
# Should-Stop:       $local_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Start inspircd
# Description:       Starts the inspircd irc server
### END INIT INFO

NAME="inspircd"
ANOPE="/usr/sbin/anoperc"
ANOPEPID="/var/run/anope.pid"
USER="irc"
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/var/lib/inspircd/services/modules

# Define LSB log_* functions (status_of_proc)
. /lib/lsb/init-functions

if [ ! -x "$ANOPE" ]; then exit 0; fi

if [ -f "$ANOPEPID" ]; then
        ANOPEPIDN=`cat "$ANOPEPID" 2> /dev/null`
fi

start_anope()
{
    [ -f "$ANOPEPID" ] || ( touch "$ANOPEPID" ; chown "$USER" "$ANOPEPID" )
        export LD_LIBRARY_PATH=/var/lib/inspircd
    start-stop-daemon --start --quiet --oknodo --chuid "$USER" --pidfile "$çPID" --exec "$ANOPE" -- start
}

stop_anope()
{
        start-stop-daemon --stop --quiet --pidfile "$ANOPEPID" > /dev/null 2> /dev/null
        rm -f "$ANOPEPID"
        return 0
}

reload_anope()
{
        if [ ! -z "$ANOPEPID" ] && kill -0 $ANOPEPID 2> /dev/null; then
                kill -HUP $ANOPEPID >/dev/null 2>&1 || return 1
                return 0
        else
                echo "Error: Anope is not running."
                return 1
        fi
}

case "$1" in
  start)
        echo -n "Starting Anope Services... "
        start_anope && echo "done."
        ;;
  stop)
        echo -n "Stopping Anope Services... "
        stop_anope && echo "done."
        ;;
  status)
    status_of_proc "$IRCD" "$NAME" && exit 0 || exit $?
    ;;
  force-reload|reload)
        echo -n "Reloading Anope Services... "
        reload_anope && echo "done."
        ;;
  restart)
        $0 stop
        sleep 2s
        $0 start
        ;;
  cron)
        start_anope || echo "Anope Services not running, starting it"
        ;;

  *)
        echo "Usage: $0 {start|stop|status|restart|reload|force-reload|cron}"
        exit 1
esac